# deploy/systemd/ai-lanbot-update.path.in
[Unit]
Description=Watch for ai-lanbot update requests

[Path]
PathModified=@INSTALL_DIR@/docker/data/update-request/request.json
Unit=ai-lanbot-update.service

[Install]
WantedBy=multi-user.target
